# monero_mining.service
[Unit]
Description=Monero mining service

[Service]
ExecStart=~/xmr-stak-cpu/bin/xmr-stak-cpu ~/xmr-stak-cpu/bin/config.txt
